# release/linux/fpm.before-remove.sh
systemctl stop chargerwalletd.service
systemctl disable chargerwalletd.service
