# release/linux/fpm.after-install.sh
systemctl enable chargerwalletd.service
systemctl start chargerwalletd.service
